# feeds/luci/applications/luci-dnsfilter/root/etc/uci-defaults/luci-dnsfilter
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@dnsfilter[-1]
	add ucitrack dnsfilter
	set ucitrack.@dnsfilter[-1].init=dnsfilter
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
